# cron.daily/remove-log
#!/bin/bash

rm -r /home/adamdrira/linkartsMJ/log;









